# cron/crontab
0 0 * * * ~/Projects/ComputerSettings/cron/cronic ~/Projects/ComputerSettings/cron/commit.sh
