# cron/crontab
@midnight sh ~/Projects/ComputerSettings/cron/commit.sh
