# cron/crontab
@midnight ~/Projects/ComputerSettings/cron/crontab sh ~/Projects/ComputerSettings/cron/commit.sh
